# Section5/dummy.service
[Unit]
Description=Our very own service
[Service]
ExecStart=/usr/sbin/dummy.sh
[Install]
WantedBy=multi-user.target
